# supervisor/conf.d/llmcli_serve.conf
[program:llmcli_serve]
command=%(ENV_HOME)s/projects/llmCLI/supervisor/scripts/run_serve.sh
directory=%(ENV_HOME)s/projects/llmCLI
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/llmCLI/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=false
autorestart=true
startsecs=20
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/llmcli/logs/llmcli_serve.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/llmcli/logs/llmcli_serve_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
