# deploy/quadlet/llmcli-fw-forwarder.container
[Unit]
Description=llmCLI Fireworks forwarder (:18646)
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=ghcr.io/roxabi/llmcli:staging
Label=io.containers.autoupdate=registry
ContainerName=llmcli-fw-forwarder
# No PublishPort — internal to roxabi.network only (reachable as
# http://llmcli-fw-forwarder:18646 from other containers on roxabi.network)
Network=roxabi.network
EnvironmentFile=%h/.roxabi/llmcli/env/proxy.env
Environment=LLMCLI_FORWARDER_PROVIDER=fireworks
Environment=LLMCLI_FORWARDER_PORT=18646
UserNS=keep-id:uid=1502,gid=1502
Exec=forwarder
HealthCmd=python -c "import urllib.request; urllib.request.urlopen('http://localhost:18646/health', timeout=5)"
HealthInterval=30s
HealthStartPeriod=15s
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=10
TimeoutStopSec=20s

[Install]
WantedBy=default.target
