# deploy/quadlet/llmcli-xai-forwarder.container
[Unit]
Description=llmCLI xAI OAuth forwarder (:18645)
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=ghcr.io/roxabi/llmcli:staging
Label=io.containers.autoupdate=registry
ContainerName=llmcli-xai-forwarder
# Reachable two ways:
#  - container-to-container via roxabi.network DNS (http://llmcli-xai-forwarder:18645)
#    — used by the LiteLLM proxy container.
#  - host loopback via PublishPort below — used by host processes such as the
#    xai-research skill (xai_research.py -> http://127.0.0.1:18645). Bound to
#    127.0.0.1 only: no LAN exposure of the OAuth-injecting relay.
PublishPort=127.0.0.1:18645:18645
Network=roxabi.network
Volume=%h/.roxabi/llmcli/credentials:/home/llmcli/.roxabi/llmcli/credentials:rw,Z
Environment=LLMCLI_FORWARDER_PROVIDER=xai
Environment=LLMCLI_FORWARDER_PORT=18645
UserNS=keep-id:uid=1502,gid=1502
Exec=forwarder
# Single-quote outer / double-quote inner: podman 5.7.0 (quadlet, M₁ Ubuntu)
# truncates a HealthCmd value that ENDS in an escaped double-quote, producing
# `sh: Unterminated quoted string` -> permanent unhealthy. Ending in the outer
# single quote sidesteps the bug; verified via `quadlet -dryrun` on 5.7.0 + 5.8.2.
HealthCmd=python -c 'import urllib.request; urllib.request.urlopen("http://localhost:18645/health", timeout=5)'
HealthInterval=30s
HealthStartPeriod=15s
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=10
TimeoutStopSec=20s
# Memory caps — contain a runaway to its cgroup instead of a global OOM freeze
# (2026-06-01 M₂ freeze). Tiny OAuth relay; 1G/1.5G is generous.
MemoryHigh=1G
MemoryMax=1536M

[Install]
WantedBy=default.target
