# deploy/quadlet/llmcli-fw-forwarder.container
[Unit]
Description=llmCLI Fireworks forwarder (:18646)
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=ghcr.io/roxabi/llmcli:staging
Label=io.containers.autoupdate=registry
ContainerName=llmcli-fw-forwarder
# No PublishPort — internal to roxabi.network only (reachable as
# http://llmcli-fw-forwarder:18646 from other containers on roxabi.network)
Network=roxabi.network
EnvironmentFile=%h/.roxabi/llmcli/env/proxy.env
Environment=LLMCLI_FORWARDER_PROVIDER=fireworks
Environment=LLMCLI_FORWARDER_PORT=18646
UserNS=keep-id:uid=1502,gid=1502
Exec=forwarder
# NB: Quadlet strips a trailing quote from a HealthCmd value that ends in one
# (python -c "...'...'..." → unterminated shell string → probe always errors →
# false unhealthy, ~13d). Ending the value with a non-quote (|| exit 1) keeps the
# inner quoting intact so this HTTP /health probe generates and runs correctly.
HealthCmd=python -c 'import urllib.request; urllib.request.urlopen("http://localhost:18646/health", timeout=5)' || exit 1
HealthInterval=30s
HealthStartPeriod=15s
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=10
TimeoutStopSec=20s
# Memory caps — contain a runaway to its cgroup instead of a global OOM freeze
# (2026-06-01 M₂ freeze). Tiny relay; 1G/1.5G is generous.
MemoryHigh=1G
MemoryMax=1536M

[Install]
WantedBy=default.target
