# deploy/quadlet/llmcli-nats-worker.container
[Unit]
Description=llmCLI NATS worker (canonical lyra.llm.generate.request)
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=ghcr.io/roxabi/llmcli:staging
Label=io.containers.autoupdate=registry
ContainerName=llmcli-nats-worker
# host network — rootless podman cannot route tailnet 100.x without it
Network=host
Volume=%h/.cache/huggingface:/home/llmcli/.cache/huggingface
Environment=HF_HOME=/home/llmcli/.cache/huggingface
# Optional: bind-mount host daemon socket to enable SWAP/STATUS pre-check on startup.
# Requires llmcli daemon running on the worker host at ~/.local/state/llmcli/llmcli.sock.
# Without this, _ensure_model auto-skips (daemon-optional mode) — safe for pure remote-worker.
# Volume=%h/.local/state/llmcli/llmcli.sock:/home/llmcli/.local/state/llmcli/llmcli.sock
Secret=llmcli-nats-nkey,type=mount,target=/home/llmcli/.config/llmcli/nkeys/seed.seed,mode=0400,uid=1502,gid=1502
Secret=llmcli-litellm-key,type=env,target=LLMCLI_LITELLM_API_KEY
AddDevice=nvidia.com/gpu=all
UserNS=keep-id:uid=1502,gid=1502
# Path namespaces: %h/.roxabi/llmcli/ (line above) is HOST-side per Roxabi data-dir
# convention. The /home/llmcli/.config/llmcli/ paths below are CONTAINER-internal
# (the llmcli service-account home in the image) where Podman mounts the secret —
# they intentionally remain on the container's XDG default and are unrelated to the
# host catalog path migration.
EnvironmentFile=%h/.roxabi/llmcli/worker.env
Environment=LLMCLI_NATS_NKEY_PATH=/home/llmcli/.config/llmcli/nkeys/seed.seed
Environment=NATS_NKEY_SEED_PATH=/home/llmcli/.config/llmcli/nkeys/seed.seed
Environment=LLMCLI_LITELLM_URL=http://localhost:4000/v1
Exec=llm --model qwen3-8b
HealthCmd=pgrep -f "llmcli nats-serve"
HealthInterval=30s
HealthStartPeriod=30s
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
